# files/etc/uci-defaults/99-auto-resize
#!/bin/sh

LOGFILE="/root/auto-resize.log"
echo "===== EXT4 终极自动扩容开始 $(date) =====" > $LOGFILE

# 1. 获取根设备 (例如 /dev/sda2)
ROOT_DEV=$(findmnt -n -o SOURCE /)
[ -z "$ROOT_DEV" ] && exit 0

# 2. 解析磁盘名和分区号
case "$ROOT_DEV" in
    *nvme*|*mmcblk*)
        DISK=$(echo $ROOT_DEV | sed 's/p[0-9]*$//')
        PART=$(echo $ROOT_DEV | sed 's/.*p//')
        ;;
    *)
        DISK=$(echo $ROOT_DEV | sed 's/[0-9]*$//')
        PART=$(echo $ROOT_DEV | sed 's/.*[^0-9]//')
        ;;
esac

# 3. 获取分区的精确起始扇区 (关键：重建必须完全对齐)
START_SECTOR=$(cat /sys/class/block/$(basename $ROOT_DEV)/start)
echo "Root: $ROOT_DEV | Disk: $DISK | Part: $PART | Start: $START_SECTOR" >> $LOGFILE

# 4. 使用 fdisk 强制重写分区表
# 逻辑：d(删除2) -> d(删除3,清障) -> n(新建) -> p(主分区) -> 2(分区号) -> 起始扇区 -> 默认终点 -> N(不抹除签名) -> w(保存)
echo "正在强行重写分区表并清空后续分区..." >> $LOGFILE
fdisk "$DISK" <<EOF >> $LOGFILE 2>&1
d
$PART
d
$((PART + 1))
d
$((PART + 2))
n
p
$PART
$START_SECTOR

N
w
EOF

# 5. 强制刷新内核分区状态
echo "刷新内核分区状态..." >> $LOGFILE
partx -u "$ROOT_DEV" >> $LOGFILE 2>&1
partprobe "$DISK" >> $LOGFILE 2>&1

# 6. 执行在线文件系统扩容
echo "执行 resize2fs..." >> $LOGFILE
resize2fs "$ROOT_DEV" >> $LOGFILE 2>&1

echo "===== 扩容任务结束 =====" >> $LOGFILE
df -h / >> $LOGFILE
