# files/etc/uci-defaults/99-auto-resize
#!/bin/sh

LOGFILE="/root/auto-resize.log"
echo "===== EXT4 GPT 专家版扩容开始 $(date) =====" > $LOGFILE

# 1. 定位根设备
ROOT_DEV=$(findmnt -n -o SOURCE /)
[ -z "$ROOT_DEV" ] && exit 0

case "$ROOT_DEV" in
    *nvme*|*mmcblk*) DISK=$(echo $ROOT_DEV | sed 's/p[0-9]*$//'); PART=$(echo $ROOT_DEV | sed 's/.*p//') ;;
    *) DISK=$(echo $ROOT_DEV | sed 's/[0-9]*$//'); PART=$(echo $ROOT_DEV | sed 's/.*[^0-9]//') ;;
esac

echo "Step 1: 修复 GPT 报头并移动备份表到磁盘末尾..." >> $LOGFILE
# -e (move backup data structures to the end of the disk) 是核心
sgdisk -e "$DISK" >> $LOGFILE 2>&1

echo "Step 2: 尝试删除阻挡分区 (如有 sda3)..." >> $LOGFILE
sgdisk -d $((PART + 1)) "$DISK" >> $LOGFILE 2>&1

echo "Step 3: 将分区 $PART 拉伸至最大剩余空间..." >> $LOGFILE
# 获取分区 $PART 的起始位置，并在原位扩容
# 使用 sgdisk 的删除再重建逻辑，或直接使用 -d / -n 但 sgdisk 有更简单的方案
# 我们可以直接删除后重新分配，只要 ID 和起点对
START_SECTOR=$(cat /sys/class/block/$(basename $ROOT_DEV)/start)
GUID=$(sgdisk -i $PART "$DISK" | grep "Partition GUID code" | awk '{print $4}')

sgdisk -d $PART "$DISK" >> $LOGFILE 2>&1
sgdisk -n $PART:$START_SECTOR:0 -t $PART:$GUID "$DISK" >> $LOGFILE 2>&1

echo "Step 4: 刷新内核并执行扩容..." >> $LOGFILE
partprobe "$DISK" >> $LOGFILE 2>&1
sync && sleep 2
resize2fs "$ROOT_DEV" >> $LOGFILE 2>&1

echo "===== 扩容任务结束 =====" >> $LOGFILE
df -h / >> $LOGFILE
